# meta-rollout/recipes-rollout/rollout/files/rollout.service
[Service]
Type=simple
Restart=always
ExecStart="python3" "/rollout-agent/agent.py"
WorkingDirectory=/rollout-agent

[Install]
WantedBy=multi-user.target
